# dist/portcullis@.service
[Unit]
Description=Portcullis proxy, service %i

[Service]
Type=notify

WorkingDirectory=%h/portcullis/build
ExecStart=%h/portcullis/build/portcullis %t/portcullis/%i.config.py
ExecReload=/bin/kill -USR1 $MAINPID
Restart=on-failure
